# rpi/root-pi/crontab
*/5 * * * * /home/pi/linux-setup/rpi/scripts/fix_wg_server.sh
0 3 * * 0 /home/pi/linux-setup/rpi/scripts/update_system.sh
